# ib.timer
[Unit]
Description=Timer for IB

[Timer]
OnCalendar=*-*-* 0,5,10,15,20:00:00

[Install]
WantedBy=timers.target
